# EE2X远航版独立项目/03-更新后端/游戏更新后端-3014/systemd/ee2x-update-yh.service
[Unit]
Description=EE2X Update YH Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ee2/ee2x-new-update/game-backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/ee2/ee2x-new-update/game-backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 3014
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
