# 更新器/update_backend_mg/systemd/ee2x-update-mg.service
[Unit]
Description=EE2X Update MG Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/ee2x/ee2x_up-mg
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/ee2x/ee2x_up-mg/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 3010
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
